# init/launchd/com.cliphub.hub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.cliphub.hub</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/cliphub</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/cliphub.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/cliphub.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CLIPHUB_MAX_HISTORY</key>
		<string>50</string>
		<key>CLIPHUB_TTL</key>
		<string>24h</string>
	</dict>
	<!-- Optional: set CLIPHUB_HOSTNAME to expose the hub under a different tailnet hostname. -->
</dict>
</plist>
